# discord-news-bot.service
[Unit]
Description=Discord AI NEWS Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=murakami
WorkingDirectory=/home/murakami/multi-agent-shogun
ExecStart=/home/murakami/multi-agent-shogun/.venv/bin/python3 scripts/discord_news_bot.py
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
